# profile.d/chocoUpgrade.alias.ps1
# upgrade
function chocoUpgrade() {
    choco upgrade all
}

Set-Alias -name "upgrade" -value "chocoUpgrade"