# profile.d/explorerHere.alias.ps1
# e.
function explorerHere() {
    Start-Process -FilePath explorer.exe -argumentlist .
}

Set-Alias -name "e." -value "explorerHere"