# profile.d/goToFolder.config.example
alias,folder
c,~\Code
d,D:\MyFolder\Downloads
p,~\Pictures